# voyanti_ems/rootfs/etc/services.d/example/run
#!/usr/bin/with-contenv bashio
# ==============================================================================
# Start the example service
# s6-overlay docs: https://github.com/just-containers/s6-overlay
# ==============================================================================

# Add your code here

# Fetch options using bashio
SITE_ID=$(bashio::config 'VOYANTI_SITE_ID')
MQTT_HOST=$(bashio::config 'MQTT_HOST')
MQTT_PORT=$(bashio::config 'MQTT_PORT')
MQTT_USERNAME=$(bashio::config 'MQTT_USERNAME')
MQTT_PASSWORD=$(bashio::config 'MQTT_PASSWORD')

# Export them as environment variables for the Python script
export VOYANTI_SITE_ID=$SITE_ID
export MQTT_HOST=$MQTT_HOST
export MQTT_PORT=$MQTT_PORT
export MQTT_USERNAME=$MQTT_USERNAME
export MQTT_PASSWORD=$MQTT_PASSWORD

## Run your program
# exec /usr/bin/my_program
exec python3 /usr/src/app.py
